# systemd/dicomq-clean.service
[Unit]
Description=dicomq orphan reaper

[Service]
Type=oneshot
ExecStart=/usr/local/bin/dicomq-clean
Environment=DICOMQ_SPOOL=/var/spool/dicomq
User=dicomq
Group=dicomq
ReadWritePaths=/var/spool/dicomq
ProtectSystem=strict
PrivateTmp=yes
NoNewPrivileges=yes
