# systemd/dicomq-send.service
[Unit]
Description=dicomq queue runner
After=network.target
RequiresMountsFor=/var/spool/dicomq

[Service]
ExecStart=/usr/local/bin/dicomq-send
Restart=always
RestartSec=5
Environment=DICOMQ_SPOOL=/var/spool/dicomq
User=dicomq
Group=dicomq
# delivering to a maildir outside the spool? list its path here too, or
# ProtectSystem=strict defers that delivery forever
ReadWritePaths=/var/spool/dicomq
ProtectSystem=strict
PrivateTmp=yes
NoNewPrivileges=yes

[Install]
WantedBy=multi-user.target
